# docker/conf/geocloud2-metrics.service
[Unit]
Description=GeoCloud2 Prometheus Metrics Server
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/geocloud2
ExecStart=/usr/bin/php /var/www/geocloud2/app/scripts/metrics-server.php
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
